# scripts/com.mohanshen.join-quant-autostudy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: resume the auto-STUDY batch on a schedule (sibling of the autoresearch one).

  Fires scripts/autostudy-loop.sh every 1 hour. Each fire resumes the pinned study session so
  the batch grinds through study/manifest.json (all normalized strategies) across quota resets,
  until every strategy is `done` or you stop. A shared jq-pipeline lock ensures study and
  research never run backtests at the same time. Chrome must stay logged in at localhost:9225.

  USAGE_LIMIT = 240 (60 free + ~180 credit-paid JQ min/day) for throughput on the large
  ~82-strategy batch. Lower it to 55 to stay inside the free tier (slower, no credit spend).

  Install:  cp scripts/com.mohanshen.join-quant-autostudy.plist ~/Library/LaunchAgents/
            launchctl load ~/Library/LaunchAgents/com.mohanshen.join-quant-autostudy.plist
  Stop:     launchctl unload ~/Library/LaunchAgents/com.mohanshen.join-quant-autostudy.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mohanshen.join-quant-autostudy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mshen/repos/join-quant/scripts/autostudy-loop.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mshen/repos/join-quant</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/mshen/.local/bin:/Users/mshen/.nvm/versions/node/v24.13.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>USAGE_LIMIT</key>
    <string>240</string>
    <key>JQ_CDP_URL</key>
    <string>http://localhost:9225</string>
  </dict>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/mshen/repos/join-quant/data/autostudy-logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mshen/repos/join-quant/data/autostudy-logs/launchd.err.log</string>
</dict>
</plist>
