# scripts/com.mohanshen.join-quant-autoresearch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: resume the join-quant autoresearch loop on a schedule.

  Fires scripts/autoresearch-loop.sh every 1 hour (StartInterval 3600s). Each fire is a
  cheap poll that exits immediately unless the epoch is started AND both the JQ daily
  backtest budget and the Anthropic plan quota are available. Polling hourly (rather than
  every 4h = the quota-window length) means a run stopped by the Anthropic per-4h-slot
  limit is picked back up within ~1h of the quota reset instead of wasting most of a window.
  A fire while still rate-limited is a no-op (claude -p exits immediately; when you're over
  the limit there is no quota to consume). Chrome must stay logged in at localhost:9225.

  Install:
    cp scripts/com.mohanshen.join-quant-autoresearch.plist ~/Library/LaunchAgents/
    launchctl load  ~/Library/LaunchAgents/com.mohanshen.join-quant-autoresearch.plist
  Verify / control:
    launchctl list | grep join-quant-autoresearch
    launchctl start com.mohanshen.join-quant-autoresearch   # fire once now
  Stop it (this is the "unless I ask it to stop" off-switch):
    launchctl unload ~/Library/LaunchAgents/com.mohanshen.join-quant-autoresearch.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mohanshen.join-quant-autoresearch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mshen/repos/join-quant/scripts/autoresearch-loop.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mshen/repos/join-quant</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/mshen/.local/bin:/Users/mshen/.nvm/versions/node/v24.13.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>USAGE_LIMIT</key>
    <string>55</string>
    <key>JQ_CDP_URL</key>
    <string>http://localhost:9225</string>
  </dict>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mshen/repos/join-quant/data/autoresearch-logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mshen/repos/join-quant/data/autoresearch-logs/launchd.err.log</string>
</dict>
</plist>
